# deploy/space-duel.service
# systemd unit for the space-duel (The Last Meridian) Colyseus game server.
#
# Install (docs/DEPLOY.md has the full provisioning checklist):
#   sudo useradd --system --home /opt/space-duel --shell /usr/sbin/nologin spaceduel
#   sudo mkdir -p /opt/space-duel && sudo chown spaceduel:spaceduel /opt/space-duel
#   # copy server/dist/server.mjs (built by `npm run build -w @space-duel/server`)
#   # to /opt/space-duel/server.mjs, then:
#   sudo cp deploy/space-duel.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now space-duel
#
# The process listens on localhost:2567 only — the TLS reverse proxy
# (deploy/Caddyfile or deploy/nginx-play.conf) is the public face.

[Unit]
Description=space-duel Colyseus game server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=spaceduel
Group=spaceduel
WorkingDirectory=/opt/space-duel
Environment=NODE_ENV=production
Environment=PORT=2567
# Cap V8's old space well under the droplet's 1GB (shared with Caddy + OS) so
# major GC runs while the heap is still resident in RAM. Uncapped, the heap
# grows until the OS swaps part of it out, and a major GC touching swapped
# pages stalls the single-threaded sim loop for SECONDS — the periodic
# freeze-then-snap every client sees (docs/perf-freeze-investigation.md,
# hypothesis 1). To verify after a deploy: add --trace-gc here temporarily and
# confirm Mark-sweep pauses stay under ~150ms in `journalctl -u space-duel`.
# NOTE: editing this file does NOT touch the droplet — reinstall the unit:
#   sudo cp deploy/space-duel.service /etc/systemd/system/ &&
#   sudo systemctl daemon-reload && sudo systemctl restart space-duel
Environment=NODE_OPTIONS=--max-old-space-size=512
ExecStart=/usr/bin/env node /opt/space-duel/server.mjs
Restart=always
RestartSec=3

# Hardening — the server needs nothing but its bundle and a listen socket.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
